# anbox-session-manager.service
[Unit]
Description=Anbox Session Manager

[Service]
ExecStart=/snap/bin/anbox session-manager

[Install]
WantedBy=default.target
