# user/weston.service
[Unit]
Description=Weston Compositor

[Service]
ExecStart=/usr/bin/weston --socket=waydroid

[Install]
WantedBy=graphical-session.target
